# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/forked-daapd.list
/usr/lib/forked-daapd/forked-daapd-sqlext.la
/etc/hotplug.d/block/20-forked-daapd-hotplug
/usr/sbin/forked-daapd
/usr/lib/forked-daapd/forked-daapd-sqlext.so
/etc/forked-daapd.conf
/etc/init.d/forked-daapd
/etc/config/forked-daapd
