# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/igmpproxy.list
/etc/init.d/igmpproxy
/usr/sbin/igmpproxy
/etc/hotplug.d/iface/40-igmpproxy
/etc/config/igmpproxy
