# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/lvm2.list
/sbin/lvmdiskscan
/sbin/vgcfgrestore
/sbin/pvs
/sbin/pvremove
/sbin/lvmsar
/sbin/vgexport
/sbin/lvrename
/sbin/lvreduce
/sbin/pvmove
/sbin/vgcfgbackup
/etc/init.d/lvm2
/sbin/vgcreate
/sbin/lvconvert
/sbin/lvremove
/sbin/vgmerge
/sbin/pvscan
/sbin/lvmconfig
/sbin/vgscan
/sbin/lvextend
/sbin/vgimport
/sbin/lvdisplay
/sbin/pvdisplay
/sbin/vgconvert
/sbin/pvresize
/sbin/pvchange
/sbin/vgck
/etc/lvm/lvm.conf
/sbin/vgmknodes
/sbin/vgreduce
/sbin/vgrename
/sbin/lvresize
/sbin/lvchange
/sbin/vgremove
/sbin/vgdisplay
/sbin/lvmchange
/sbin/lvscan
/sbin/pvcreate
/sbin/vgextend
/sbin/vgsplit
/sbin/lvm
/sbin/lvs
/sbin/lvcreate
/sbin/pvck
/sbin/lvmsadc
/sbin/vgs
/sbin/vgchange
