# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/minidlna.list
/etc/init.d/minidlna
/usr/bin/minidlna
/etc/hotplug.d/block/20-minidlna-hotplug
/etc/config/minidlna
