# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/netatalk.list
/usr/sbin/netatalk
/usr/lib/uams/uams_clrtxt.so
/usr/lib/uams/uams_pam.so
/usr/lib/uams/uams_dhx.so
/usr/lib/uams/uams_passwd.so
/usr/lib/uams/uams_dhx_pam.so
/usr/bin/dbd
/usr/lib/uams/uams_dhx2_pam.so
/usr/sbin/cnid_metad
/usr/lib/libatalk.a
/usr/lib/uams/uams_randnum.so
/usr/bin/afppasswd
/usr/lib/uams/uams_dhx2_passwd.so
/usr/lib/libatalk.so.16
/usr/lib/libatalk.la
/usr/bin/apple_dump
/usr/bin/ad
/usr/lib/uams/uams_guest.so
/usr/lib/uams/uams_dhx_passwd.so
/etc/init.d/netatalk
/usr/lib/libatalk.so
/usr/sbin/cnid_dbd
/usr/bin/cnid2_create
/usr/lib/uams/uams_dhx2.so
/etc/afp.conf
/usr/lib/libatalk.so.16.0.0
/etc/config/netatalk
/usr/sbin/afpd
