# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/rpcd.list
/etc/init.d/rpcd
/etc/config/rpcd
/usr/share/rpcd/acl.d/superuser.json
/sbin/rpcd
/usr/share/rpcd/acl.d/unauthenticated.json
