# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/samba4-server.list
/usr/bin/smbcontrol
/etc/config/samba
/usr/bin/smbstatus
/usr/bin/smbtree
/usr/sbin/nmbd
/usr/bin/smbcacls
/usr/bin/testparm
/etc/samba/smb.conf.template
/etc/init.d/samba
/usr/bin/smbcquotas
/usr/sbin/pdbedit
/usr/sbin/smbpasswd
/lib/upgrade/keep.d/samba4-server
/usr/sbin/smbd
/usr/bin/nmblookup
/usr/sbin/samba4_multicall
