# OpenWrt-ImageBuilder-rtd1295-mnas_emmc.Linux-x86_64/target/linux/rtd1295/image/image_file-r160868/tmp/rootfs/usr/lib/opkg/info/vsftpd-tls.list
/etc/config/vsftpd
/lib/upgrade/keep.d/vsftpd-tls
/usr/sbin/vsftpd
/etc/init.d/vsftpd
